# backend/rs-backend.service
[Unit]
Description=Router Scan Backend Service
After=network.target

[Service]
ExecStart=/home/kali/router-scan-command/backend/rs-backend.sh
Type=oneshot
RemainAfterExit=true

[Install]
WantedBy=default.target