# backend/rs-backend.service
[Unit]
Description=My React App Service
After=network.target

[Service]
User=root
WorkingDirectory=/home/kali/router-scan-backend
ExecStart=/home/kali/router-scan-command/backend.service.sh

[Install]
WantedBy=multi-user.target